# rootfs/usr/lib/systemd/system/first-boot-setup.service
[Unit]
Description=First Boot Setup: Add Flathub Repositories, Install Applications, and Apply Overrides
ConditionPathExists=!/var/lib/flatpak/.first-boot-setup-done
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/gnome-terminal -- bash -c "/usr/local/bin/first-boot-setup.sh; exec bash"
ExecStartPost=/usr/bin/touch /var/lib/flatpak/.first-boot-setup-done

[Install]
WantedBy=multi-user.target



